# systemd/darkwall-comfyui.service
[Unit]
Description=Generate dark-mode wallpaper using ComfyUI
Documentation=man:darkwall(1)

# Only run when network is available (for ComfyUI communication)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.nix-profile/bin/darkwall generate

# Environment variables (optional)
Environment=DARKWALL_LOG_LEVEL=INFO

# Security settings for user service
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
TimeoutStartSec=300
TimeoutStopSec=30

# Restart on failure
Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
