# systemd/darkwall-comfyui.service
[Unit]
Description=Generate dark-mode wallpaper using ComfyUI
Documentation=man:generate-wallpaper-once(1)

# Only run when network is available (for ComfyUI communication)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.nix-profile/bin/generate-wallpaper-once generate

# Environment variables (optional - configure as needed)
Environment=DARKWALL_LOG_LEVEL=INFO
#Environment=COMFYUI_BASE_URL=https://comfyui.home.arpa
#Environment=TIME_SLOT_MINUTES=30
#Environment=WALLPAPER_OUTPUT_PATH=%h/Pictures/wallpapers/current.png

# Security settings for user service
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
TimeoutStartSec=300
TimeoutStopSec=30

# Restart on failure
Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
